# deploy/com.linkqueue.triage.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.linkqueue.triage</string>
    <key>ProgramArguments</key>
    <array>
        <!-- uv tool install puts the binary here; adjust for your username -->
        <string>/Users/sagarpal/.local/bin/obs_triage</string>
        <string>run</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/linkqueue-triage.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/linkqueue-triage.log</string>
</dict>
</plist>
